# 0x1A-application_server/gunicorn.service
[Unit]
Description=Gunicorn instance for HBNB
After=network.target

[Service]
PermissionsStartOnly=true
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/AirBnB_clone_v4

# Ensure /run/hbnb exists with correct ownership
RuntimeDirectory=hbnb
RuntimeDirectoryMode=755

# Set environment variables
Environment="HBNB_MYSQL_USER=hbnb_dev"
Environment="HBNB_MYSQL_PWD=hbnb_dev_pwd"
Environment="HBNB_MYSQL_HOST=localhost"
Environment="HBNB_MYSQL_DB=hbnb_dev_db"
Environment="HBNB_TYPE_STORAGE=db"

# Start Gunicorn with 4+ workers
ExecStart=/home/ubuntu/AirBnB_clone_v4/venv/bin/gunicorn --workers 4 \
    --access-logfile /tmp/airbnb-access.log \
    --error-logfile /tmp/airbnb-error.log \
    --bind 0.0.0.0:5003 web_dynamic.2-hbnb:app

ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s TERM $MAINPID

PrivateTmp=false

[Install]
WantedBy=multi-user.target
